# com.yanbinwang.paper-tracker.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.yanbinwang.paper-tracker</string>
  <key>ProgramArguments</key><array>
    <string>/bin/sh</string>
    <string>REPLACE_WITH_ABSOLUTE_PATH/sync.sh</string>
  </array>
  <key>StartInterval</key><integer>7200</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>REPLACE_WITH_ABSOLUTE_PATH/tracker.log</string>
  <key>StandardErrorPath</key><string>REPLACE_WITH_ABSOLUTE_PATH/tracker.log</string>
</dict></plist>

